# com.ztrade.marketopen.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ztrade.marketopen</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/aboros/Ztrade/start_trading_at_market_open.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/aboros/Ztrade</string>

    <!-- Run at 9:25 AM ET Monday-Friday -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>25</integer>
        </dict>
        <!-- Tuesday -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>25</integer>
        </dict>
        <!-- Wednesday -->
        <dict>
            <key>Weekday</key>
            <integer>3</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>25</integer>
        </dict>
        <!-- Thursday -->
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>25</integer>
        </dict>
        <!-- Friday -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>25</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/aboros/Ztrade/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/aboros/Ztrade/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Don't run if battery is low -->
    <key>LaunchOnlyOnce</key>
    <false/>

    <!-- Run in user context -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
